# netrunner/deploy/netrunner.service
[Unit]
Description=NETRUNNER Backend
After=network.target

[Service]
User=netrunner
WorkingDirectory=/opt/netrunner/backend
Environment=NETRUNNER_PASSWORD=DEIN_PASSWORT_HIER
Environment=NETRUNNER_JWT_SECRET=DEIN_JWT_SECRET_HIER
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
